# config/vx-election-manager.service
[Unit]
Description=VotingWorks Election Manager

[Service]
Type=simple
User=vx-services
Environment=VX_CONFIG_ROOT=/vx/config
Environment=MODULE_SEMS_CONVERTER_WORKSPACE=/vx/data/module-sems-converter
ExecStart=/bin/bash /vx/services/run-election-manager.sh
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=election-manager

[Install]
WantedBy=multi-user.target