# config/vx-festival.service
[Unit]
Description=Festival Speech Synthesis

[Service]
Type=simple
ExecStart=/usr/bin/festival --server

[Install]
WantedBy=multi-user.target